# docker/fs/etc/service/wireport-server/finish
#!/bin/sh

echo "> Stopping wireport-server service"

pkill -TERM wireport
