# docker/fs/etc/service/wireport-server/run
#!/bin/sh

exec /usr/bin/wireport server start
