# docker/fs/etc/service/wireport-server/run
#!/bin/sh

echo "> Starting wireport-server service"

exec /usr/bin/wireport server start
