# docker/fs/etc/service/wireguard/finish
#!/bin/sh

echo "> Stopping wireguard service"

wg-quick down wg0
