# docker/fs/etc/service/caddy/finish
#!/bin/sh

echo "> Stopping caddy service"

pkill -TERM caddy
