# docker/fs/etc/service/caddy/run
#!/bin/sh

exec /usr/bin/caddy run --config /etc/caddy/Caddyfile
