# docker/fs/etc/service/caddy/run
#!/bin/sh

exec env GODEBUG=netdns=cgo /usr/bin/caddy run --config /etc/caddy/Caddyfile
