# docker/fs/etc/service/caddy/run
#!/bin/sh

exec env GODEBUG=netdns=cgo XDG_DATA_HOME=/data XDG_CONFIG_HOME=/config /usr/bin/caddy run --config /etc/caddy/Caddyfile
