# docker/fs/etc/service/coredns/finish
#!/bin/sh

echo "> Stopping coredns service"

pkill -TERM coredns
