# docker/fs/etc/service/coredns/run
#!/bin/sh

exec /usr/bin/coredns -conf /etc/coredns/Corefile
