# docker/fs/etc/service/coredns/run
#!/bin/sh

echo "> Starting coredns service"

exec /usr/bin/coredns -conf /etc/coredns/Corefile
