# docker/fs/etc/service/wireport/run
#!/bin/sh

exec /usr/bin/wireport host start
