# docker/fs/etc/service/wireport-gateway/finish
#!/bin/sh

echo "> Stopping wireport-gateway service"

pkill -TERM wireport
