# docker/fs/etc/service/wireport-gateway/run
#!/bin/sh

echo "> Starting wireport-gateway service"

exec /usr/bin/wireport gateway start
